# lib/systemd/system/plugh.socket
[Unit]
Description=Socket Activation plugh socket

[Socket]
ListenStream=50004
#Accept=yes

[Install]
WantedBy=sockets.target
